# deploy/agentdeck-smoke.timer
[Unit]
Description=agentdeck nightly smoke at 04:30

[Timer]
OnCalendar=*-*-* 04:30:00
Persistent=true

[Install]
WantedBy=timers.target
