# deploy/agentdeck.service
# systemd unit for bare-metal installs (recommended on the homelab).
#   sudo cp deploy/agentdeck.service /etc/systemd/system/
#   sudo systemctl enable --now agentdeck
[Unit]
Description=agentdeck — mission control for AI coding agents
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=admin
WorkingDirectory=/home/admin/projects/agentdeck
# URL that TARGETS use to reach this server for approval callbacks
Environment=AGENTDECK_BASE_URL=http://YOUR-SERVER-IP:9110
ExecStart=/home/admin/projects/agentdeck/.venv/bin/python -m server
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
