# deploy/agentdeck-smoke.service
[Unit]
Description=agentdeck nightly smoke test
After=agentdeck.service

[Service]
Type=oneshot
User=admin
# Discord bot token + channel live here (mode 600, NOT in the repo):
#   AGENTDECK_DISCORD_BOT_TOKEN=...
#   AGENTDECK_DISCORD_CHANNEL_ID=...
EnvironmentFile=-/etc/agentdeck-smoke.env
ExecStart=/usr/bin/python3 /home/admin/projects/agentdeck/deploy/agentdeck-smoke.py
TimeoutStartSec=1200
